# docker/runit/template/run
#!/bin/sh -e
# redirects all stderr output to stdout
exec 2>&1
exec COMMAND
