# docker/runit/php-fpm/run
#!/bin/sh -e
# redirects all stderr output to stdout
#exec 2>&1
exec sh -c "/usr/local/sbin/php-fpm"
